# infra/ansible/roles/app/templates/tbc-worker-brief.service.j2
[Unit]
Description=TBC Worker — Morning Brief scheduler (Sonnet 4.6 daily digest)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=tbc
Group=tbc
WorkingDirectory=/opt/tbc
EnvironmentFile=/etc/tbc/env
ExecStart=uv run python -m services.worker_brief.main
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tbc-worker-brief

[Install]
WantedBy=multi-user.target
